# Build/supervisord/fcgisock.conf
[program:fcgisock]
user = root
command = %(ENV_SPAWN_FCGI)s -F 1 -s /var/run/fcgiwrap_%(ENV_NAME)s.socket -U www-data -G www-data -n -- %(ENV_DAEMON)s -f -p /usr/lib/cgi-bin/nut/upsimage.cgi  -p /usr/lib/cgi-bin/nut/upsset.cgi  -p /usr/lib/cgi-bin/nut/upsstats.cgi 2>&1
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
autostart = true
startretries = 5
startsecs = 0
priority = 10