# Build/supervisord/fcgiwrap.conf
[program:fcgiwrap]
user = root
command = %(ENV_SPAWN_FCGI)s -a %(ENV_FCGI_ADDR)s -p %(ENV_FCGI_PORT)s -- %(ENV_DAEMON)s -f 2>&1
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
autostart = true
startretries = 0
startsecs = 0
priority = 10