# Build/supervisord/script.conf
[program:start]
command=/start.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
autostart=true
startretries=0
startsecs=0
